# vagrant_scripts/celery_upstart.conf
description "Celery queue process worker for swarfarm.com"

start on vagrant-mounted
stop on runlevel [!2345]

respawn
setuid vagrant
setgid www-data
kill timeout 120

exec /home/vagrant/.pyenv/versions/swarfarm-3.6.7/bin/celery --workdir=/vagrant -A swarfarm.celery worker --loglevel=info
